# 09_FMS/단말/01_Source/linux/crontab/del_file.sh
#!/bin/sh
SERVICE='Delete_OldFile'

find /home/fms/log -ctime +1 -exec rm {} \;
